# systemd/user/mpd.service.d/override.conf
[Service]
Environment=XDG_CONFIG_HOME=/home/gerze/.config
Environment=XDG_DATA_HOME=/home/gerze/.local/share
ExecStart=
ExecStart=/usr/bin/mpd --systemd ${XDG_CONFIG_HOME}/mpd/mpd.conf
