# systemd/user/gtklock.service
[Unit]
Description=GTKLock Screen Locker
Documentation=man:gtklock(1)

[Service]
Type=simple
ExecStart=/usr/bin/gtklock
Restart=no
